# vendor/github.com/getlantern/proxybench/http-proxy.conf
# http-proxy - upstart job file

description "lantern server-side proxy"
author "Snyder Pearson <aranhoide@getlantern.org>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Terminate with SIGINT instead of SIGKILL
kill signal SIGINT

# Wait 15 seconds before sending SIGKILL if SIGINT doesn't work
kill timeout 15

# Run as
setuid root

# Run in home directory
chdir /root

# We've observed we need to raise the limit on open file descriptors
limit nofile 1024768 1024768

# Start the process, piping output prepended with a prefix to syslog
exec /root/http-proxy -bench -addr :443 -obfs4-addr :24443 -lampshade-addr :24444 -https -key key.pem -cert cert.pem -tos 0 2>&1 | logger -t http-proxy
